# provision/post/mailhog/mailhog.conf
description "Mailhog"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
pre-start script
    exec su - vagrant -c "/usr/bin/env ~/bin/mailhog > /dev/null 2>&1 &"
end script
